# deploy/skyview.service
[Unit]
Description=Skyview Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/skyview
Environment=PYTHONPATH=/opt/skyview/backend
ExecStart=/opt/skyview/venv/bin/python /opt/skyview/backend/app.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/skyview.log
StandardError=append:/var/log/skyview.log

[Install]
WantedBy=multi-user.target
